# .docker/s6-rc.d/octane/run
#!/bin/sh -e

php /var/www/artisan octane:start
